# systemd/rag-medicine-instructions.service
[Unit]
Description=RAG Medicine Instructions Streamlit Application
After=network.target

[Service]
Type=simple
User=rag-app
Group=rag-app
WorkingDirectory=/opt/rag-medicine-instructions
Environment="PATH=/opt/rag-medicine-instructions/.venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/opt/rag-medicine-instructions/.env

# Streamlit command
ExecStart=/opt/rag-medicine-instructions/.venv/bin/streamlit run app/streamlit_app.py \
    --server.port=8501 \
    --server.address=127.0.0.1 \
    --server.headless=true \
    --server.enableCORS=false \
    --server.enableXsrfProtection=false

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rag-medicine-instructions

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/rag-medicine-instructions/storage /opt/rag-medicine-instructions/data /var/log/rag-medicine-instructions

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target

